# raspberry-pi-4/rc.local
#!/bin/bash 
/configure_te_pi.sh 
exit 0 